# docs/deploy/com.aitelier.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd template for the aitelier service supervisor.
  @REPO_ROOT@ and @PATH@ are substituted by scripts/install-launchd.sh.
  Installed to ~/Library/LaunchAgents/com.aitelier.agent.plist.

  RunAtLoad starts it at login; KeepAlive restarts it whenever it exits
  (crash or manual kill); ThrottleInterval throttles crash loops.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aitelier.agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>@REPO_ROOT@/scripts/supervise.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>@REPO_ROOT@</string>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>@PATH@</string>
    </dict>

    <key>StandardOutPath</key>
    <string>@REPO_ROOT@/runs/logs/aitelier.log</string>
    <key>StandardErrorPath</key>
    <string>@REPO_ROOT@/runs/logs/aitelier.log</string>
</dict>
</plist>
